# systemd/lianli-evdi-setup.service
[Unit]
Description=Pre-create evdi virtual display node for Lian Li daemon
ConditionPathExists=/sys/devices/evdi/add
After=systemd-modules-load.service
Before=lianli-daemon@.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo 1 > /sys/devices/evdi/add'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
